# src/var/lib/dpkg/info/calendar.list
/etc
/etc/calendar
/etc/calendar/default
/etc/cron.daily
/etc/cron.daily/calendar
/etc/default
/etc/default/calendar
